# python-webserver/pythonapi.service
[Unit]
After=network.target

[Service]
ExecStart=python3 /usr/local/bin/pythonWebServer.py 8080

[Install]
WantedBy=default.target